# .config/etc/systemd/system/container-entrypoint.target
[Unit]
Description=the target for container-entrypoint.service
Requires=container-entrypoint.service systemd-logind.service systemd-user-sessions.service

