# .config/etc/systemd/system/docker-entrypoint.target
[Unit]
Description=the target for docker-entrypoint.service
Requires=docker-entrypoint.service systemd-logind.service systemd-user-sessions.service

